# tests/fixtures/gluon/execution/source-trees/cast-daemon-fixture-1.0.0/packaging/cast-daemon.service.in
# SPDX-FileCopyrightText: 2026 AerynOS Developers
# SPDX-License-Identifier: MPL-2.0

[Unit]
Description=Cast daemon execution fixture @PROJECT_VERSION@
After=network.target

[Service]
Type=oneshot
ExecStart=@CMAKE_INSTALL_FULL_SBINDIR@/cast-daemon-fixture --check-config @CMAKE_INSTALL_FULL_DATADIR@/defaults/cast-daemon-fixture/cast-daemon.conf
DynamicUser=yes
StateDirectory=cast-daemon-fixture

[Install]
WantedBy=multi-user.target
